# applications/luci-app-gyzddns/root/etc/uci-defaults/luci-app-gyzddns
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@gyzddns[-1]
	add ucitrack gyzddns
	set ucitrack.@gyzddns[-1].init=gyzddns
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
